# components/scip-auth-ad/systemd/system/systemd-logind.service.d/00-after-dbus.conf
# https://bugzilla.redhat.com/show_bug.cgi?id=1532105
[Unit]
After=dbus.service
BindsTo=dbus.service
